# systemd/system/git-bro.service
[Unit]
Description=Monitor git repositories

[Service]
Environment="SYSTEM_MODE=true"
ExecStart=%h/.local/bin/git-bro
