# systemd/system/git-bro.service
[Unit]
Description=Monitor git repositories

[Service]
Environment="SYSTEM_MODE=true"
ExecStart=/usr/bin/git-bro
