# systemd/git-bro.service
[Unit]
Description=Monitor git repositories

[Service]
ExecStart=%h/.local/bin/git-bro.rb
