# systemd/user/git-bro.service
[Unit]
Description=Monitor git repositories

[Service]
ExecStart=git-bro
